# crontab
TZ=America/Boise
13 18 * * * bash /home/schmuck/script.sh > /home/schmuck/log.log 2>&1
