# crontab
TZ=America/Boise
45 8 * * * bash /home/schmuck/script.sh > /home/schmuck/log.log 2>&1
