# crontab
TZ=America/Boise
0 10 * * * bash /home/schmuck/script.sh > /home/schmuck/log.log 2>&1
