# crontab
00 16 * * * /script.sh
#Empty line is required at the end of this file for a valid cron file.
